# macos/sh.ntfy.subscriber.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ntfy subscriber launch agent for macOS
  Copy to: ~/Library/LaunchAgents/sh.ntfy.subscriber.plist

  Replace /Users/YOUR_USERNAME with your actual home directory.

  Load:   launchctl load ~/Library/LaunchAgents/sh.ntfy.subscriber.plist
  Unload: launchctl unload ~/Library/LaunchAgents/sh.ntfy.subscriber.plist
  Status: launchctl list | grep ntfy
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>sh.ntfy.subscriber</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/ntfy</string>
        <string>subscribe</string>
        <string>--from-config</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/YOUR_USERNAME</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/ntfy-subscriber.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/ntfy-subscriber.err</string>
</dict>
</plist>
